# stages/04.configure-desktop-env/02.vnc-server/files/xserver.service
[Unit]
Description=ADI X Server
Wants=display-manager.service
After=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/adi-xserver.sh

[Install]
WantedBy=multi-user.target
